# android-service@hwcomposer.service.d/20-cutie.conf
[Unit]
BindsTo=cutie-shell.service
Before=cutie-shell.service
